# data/unity8-mir.conf.in
description "Unity Shell v8 on Mir"

emits scope-ui-starting

start on xsession SESSION=unity8-mir and started dbus
stop on desktop-end

respawn
expect stop

env QT_QPA_PLATFORM=ubuntumirserver
env UPSTART_JOB=unity8

pre-start script
    if [ -n "$MIR_SOCKET" ]; then
        if [ -z "$UNITY_MIR_SOCKET" ]; then
            # Save original value of MIR_SOCKET in case we are restarted,
            # as we modify the variable for future jobs, including ourself.
            UNITY_MIR_SOCKET=$MIR_SOCKET
            initctl set-env --global UNITY_MIR_SOCKET=$UNITY_MIR_SOCKET
        fi

        # Point unity8 at unity-system-compositor
        MIR_SERVER_FILE=$XDG_RUNTIME_DIR/mir_socket
        initctl set-env MIR_SERVER_FILE=$MIR_SERVER_FILE
        initctl set-env MIR_SERVER_HOST_SOCKET=$UNITY_MIR_SOCKET

        # Point all future jobs in this session to our Mir socket instead of
        # unity-system-compositor's socket.
        initctl set-env --global MIR_SOCKET=$MIR_SERVER_FILE
        initctl set-env --global QT_QPA_PLATFORM=ubuntumirclient
        initctl set-env --global EGL_PLATFORM=mir
        initctl set-env --global UBUNTU_PLATFORM_API_BACKEND=libubuntu_application_api_test.so.1
        initctl set-env --global UBUNTU_PLATFORM_API_SENSOR_TEST=@PKGDATA_DIR@/sensors.txt
        initctl set-env --global QML2_IMPORT_PATH=@LIB_DIR@/qt5/imports/Unity-Mir
        initctl set-env --global UNITY_INDICATOR_PROFILE=desktop

        gdbus call --session \
                   --dest org.freedesktop.DBus \
                   --object-path /org/freedesktop/DBus \
                   --method org.freedesktop.DBus.UpdateActivationEnvironment \
                   "@a{ss} {'QT_QPA_PLATFORM': 'ubuntumirclient', 'MIR_SOCKET': '$MIR_SERVER_FILE'}"
    fi

    if [ -z "$UNITY_SCOPES_LIST" ]; then
        initctl set-env UNITY_SCOPES_LIST="scopes;clickscope;musicaggregator;videoaggregator"
    fi

    initctl emit scope-ui-starting
end script

exec ${BINARY:-unity8} $ARGS

pre-stop exec initctl emit --no-wait desktop-end
post-stop exec dbus-send --type=method_call --address=$UPSTART_SESSION /com/ubuntu/Upstart com.ubuntu.Upstart0_6.EndSession
